# systemd/cofiswarm.target
[Unit]
Description=Cofiswarm stack (host services + compose profile)
After=network-online.target docker.service
Wants=network-online.target

[Install]
WantedBy=multi-user.target
